# lab-base-container/etc/supervisor/conf.d/ssh.conf
[program:ssh-server]
command=/usr/sbin/sshd -D
process_name = ssh-server
